# home/etc/systemd/system/backup.service
[Unit]
Description=Create a backup chronologically with rsync
Wants=daily.timer


[Service]
Type=oneshot
ExecStart=/bin/bash /home/matt/.backup.sh
User=matt


[Install]
WantedBy=multi-user.target
