# home/etc/systemd/system/daily.timer
[Unit]
Description=Create daily agenda
Requires=journal.service

[Timer]
OnCalendar=*-*-* 8:00:00
Persistent=true

[Install]
WantedBy=timers.target

